# src/Service/susi.service
[Unit]
Description=Central service to start all Susi daemons
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=andreas
ExecStart=/usr/local/bin/susi.start
ExecStop=/usr/local/bin/susi.stop
Environment="SUSI_INSTALLATION=/home/andreas/Documents/Projekte/2019-Susi/Susi"

[Install]
WantedBy=multi-user.target
